# trading_bot.service
[Unit]
Description=OANDA Trading Bot
After=network.target

[Service]
ExecStart=/usr/bin/python3.11 /path/to/bot.py
WorkingDirectory=/path/to/
Environment=OANDA_API_TOKEN=your_api_token_here
Environment=OANDA_ACCOUNT_ID=your_account_id_here
Environment=TELEGRAM_TOKEN=your_telegram_token_here
Environment=TELEGRAM_CHAT_ID=your_telegram_chat_id_here
Restart=always
User=your_user
Group=your_group

[Install]
WantedBy=multi-user.target
